# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:mongod]
command=mongod

[program:postgres]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf

[program:memcache]
command=memcached -m 64

[program:nginx]
command=nginx

[program:api]
command=python /app/api/rpgthing/manage.py runserver

[program:game_server]
command=/app/serverSide/bin/RPGThingServerSide --addr :8080
